# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── Telegram local Bot API server ──────────────────────────────────────────
[program:tg-api]
command=telegram-bot-api
    --api-id=%(ENV_TELEGRAM_API_ID)s
    --api-hash=%(ENV_TELEGRAM_API_HASH)s
    --local
    --dir=/tmp/tg-data
    --temp-dir=/tmp/tg-temp
    --http-port=8081
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── AMIGO Python bot ────────────────────────────────────────────────────────
; startsecs=8 gives the API server time to initialise before the bot connects
[program:amigo-bot]
command=python bot.py
directory=/app
autostart=true
autorestart=true
startsecs=8
startretries=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
